# files/systemd/aws.target
[Unit]
Description="AWS stuff is ready"
BindsTo=network.target
After=network.target

[Install]
WantedBy=ready.target
